# supervisord.conf
[program:libreoffice]
command=/usr/lib/libreoffice/program/soffice.bin --accept='socket,host=localhost,port=2002;urp;StarOffice.ServiceManager' --headless
autorestart=true
